# service/plantstation-wait.path
[Unit]
Description=Plant Station Delay Path Unit

[Path]
PathExists=/dev/i2c-1
PathExists=/var/lib/systemd/timesync/clock
