# service/plantstation-wait.service
[Unit]
Description=Plant Station Delay Service

[Service]
Type=oneshot
ExecStart=/opt/bin/wait-for-timesync.sh

[Install]
WantedBy=default.target
